# Setmana2/monitoritzar_logs.service
[Unit]
Description=Service to check errors from services.
After=network.target

[Service]
EnvironmentFile=/home/GSX/services
ExecStart=/home/GSX/monitoritzar_logs.sh
User=root
StandardOutput=append:/var/log/monitoritzar_logs.log
StandardError=append:/var/log/monitoritzar_logs.log

[Install]
WantedBy=multi-user.target
